# config/cron/crontab
* * * * * application /usr/local/bin/php /app/web/typo3/cli_dispatch.phpsh scheduler
