# config/cron/crontab
* * * * * root /usr/bin/php /app/web/typo3/cli_dispatch.phpsh scheduler
