# config/cron/crontab
MAILTO="infrastructure@ruhmesmeile.com"
* * * * * application /usr/local/bin/php -d "apc.enablie_cli=1" /app/vendor/bin/typo3cms scheduler:run
