# config/cron/crontab
MAILTO="infrastructure@ruhmesmeile.com"
* * * * * application /usr/local/bin/php /app/vendor/bin/typo3cms scheduler:run
