# config/cron/crontab
MAILTO="infrastructure@ruhmesmeile.com"
* * * * * application /usr/local/bin/php /app/web/typo3/sysext/core/bin/typo3 scheduler:run
